# etc/systemd/system/display-brightness-down.service
[Unit]
Description=Set display brightnesses to night time preferences

[Service]
ExecStart=/usr/local/bin/display-down.sh
Type=oneshot
User=root
