# etc/systemd/system/display-brightness-up.service
[Unit]
Description=Set display brightnesses to default

[Service]
ExecStart=/usr/local/bin/display-up.sh
Type=oneshot
User=root
